# crontab
* * * * * root . /root/project_env.sh; . /app/backup.sh
